# com.cyberbrief.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cyberbrief.refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/msk/cc/cyber/run_pipeline.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/msk/.local/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>StartInterval</key>
    <integer>43200</integer>

    <key>StandardOutPath</key>
    <string>/Users/msk/cc/cyber/data/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/msk/cc/cyber/data/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/msk/cc/cyber</string>
</dict>
</plist>
